# record-collection.service
[Unit]
Description=Record Collection Service
After=network.target

[Service]
ExecStart=/usr/bin/node /Users/harryshapiro/developer/record-collection/packages/server/dist/index.js
Restart=always
User=harryshapiro
Environment=NODE_ENV=production
WorkingDirectory=/Users/harryshapiro/developer/record-collection/packages/server

[Install]
WantedBy=multi-user.target